# etc/cron/crontab
0,10,20,30,40,50 * * * * /home/babarot/fukushima.sh
